# redis/redis-supervisor.conf
[program:redis-supervisor]
command=/home/redis/run-redis.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/redis/redis.err.log
stdout_logfile=/var/log/redis/redis.out.log
